# core_install/templates/stats/admin/raw/etc/systemd/system/eniq/runtime/eniq-webserver.service
[Unit]
Description=The Eniq Webserver Service
Requires=NAS-online.service
After=NAS-online.service

[Service]
Type=forking
User=dcuser
Group=dc5000
ExecStart=/eniq/smf/bin/eniq_smf_start_stop.sh -a start -s webserver
ExecStop=/eniq/smf/bin/eniq_smf_start_stop.sh -a stop -s webserver
Restart=always
RestartSec=15
StartLimitInterval=450
StartLimitBurst=3
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target NAS-online.service
