# core_install/templates/stats/admin/raw/etc/systemd/system/eniq/runtime/eniq-connectd.service
[Unit]
Description=ENIQ OSS connectd Service
Requires=NAS-online.service
After=NAS-online.service

[Service]
Type=forking
ExecStart=/eniq/smf/connectd/eniq_smf_connectd.sh -a start
ExecStop=/eniq/smf/connectd/eniq_smf_connectd.sh -a stop
TimeoutStartSec=180
TimeoutStopSec=120
StartLimitInterval=620
StartLimitBurst=3
Restart=always
RestartSec=15

[Install]
WantedBy= multi-user.target NAS-online.service
